# data/meshnet.service
[Unit]
Description=MeshNET Daemon

[Service]
Type=simple
ExecStart=/usr/local/sbin/meshnet

[Install]
WantedBy=network.target

