# data/meshnet.service.in
[Unit]
Description=MeshNET Daemon

[Service]
Type=simple
ExecStart=@sbindir@/meshnet

[Install]
WantedBy=network.target

